# scripts/com.chebe24.gitcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chebe24.gitcheck</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/caryhebert/Developer.nosync/nexus-command/scripts/check-git-in-documents.sh</string>
  </array>

  <!-- Run every 6 hours -->
  <key>StartInterval</key>
  <integer>21600</integer>

  <!-- Also run at login -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/caryhebert/Library/Logs/git-in-documents.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/caryhebert/Library/Logs/git-in-documents-error.log</string>
</dict>
</plist>
